# layers/meta-clockos/recipes-core/earsys/files/earsys.service
[Unit]
Description=EARSYS — EAR 기반 졸음 감지 시스템
Documentation=https://github.com/WhiPaper/EARSYS
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/opt/earsys

EnvironmentFile=-/etc/earsys.env
Environment="EARSYS_LOG_LEVEL=INFO"
Environment="EARSYS_EAR_THRESHOLD=0.23"
Environment="EARSYS_CLOSED_FRAMES=20"
Environment="EARSYS_MODEL_PATH=/opt/earsys/face_landmarker.task"

ExecStart=/usr/bin/python3 /opt/earsys/main.py

Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=earsys

# 비정상 종료 시 최대 5회 재시작 (10분 이내)
StartLimitInterval=600
StartLimitBurst=5

# 보안 강화
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=graphical.target
